# systemd/user/battery_alert.service
[Unit]
Description=Battery Alert Service

[Service]
ExecStart=/bin/bash ~/battery_alert.sh
