# systemd/user/battery_alert.timer
[Unit]
Description=Run Battery Alert Script Every 5 Minutes

[Timer]
OnUnitActiveSec=5m
Unit=battery_alert.service

[Install]
WantedBy=timers.target
